# doc/deploy/init_server/etc/supervisor/log_daemon.conf
[program:log_daemon]
user=app
environment=HOME="/home/app",USER="app"
command = bash -l -c 'ruby /home/app/railsapps/Fandom/current/bin/log_daemon.rb production /home/app/railsapps/Fandom/current 3'